# DailyJob/bug_watch/com.lis.bug-watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lis.bug-watch</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-s</string>
        <string>/bin/bash</string>
        <string>/Users/hung.l/src/vibrant-america-working-agent/DailyJob/bug_watch/run_bug_watch.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>11</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>17</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/hung.l/src/vibrant-america-working-agent/DailyJob/bug_watch/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/hung.l/src/vibrant-america-working-agent/DailyJob/bug_watch/launchd_stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/hung.l/src/vibrant-america-working-agent</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/hung.l</string>
    </dict>
</dict>
</plist>
