# scripts/com.lis-code-agent.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lis-code-agent.weekly</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/hung.l/src/lis-code-agent/scripts/weekly-routine.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>11</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/hung.l/src/lis-code-agent/logs/launchd-weekly-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/hung.l/src/lis-code-agent/logs/launchd-weekly-stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/hung.l/src/lis-code-agent</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
